# build/crontabs/qateam/hgwbeta.crontab
SHELL=/bin/csh
# specifying path so that the "which" command will work within qa scripts.
PATH=/cluster/bin/scripts:/cluster/bin/x86_64:/usr/local/bin:/bin:/usr/bin

# Minute Hour DOM Month DOW command

MAILTO=kuhn,rhead
# performance log
5,20,35,50 * * * * hgTracksRandom /data/perf/machines >> /data/perf/hgTracksRandom.log

MAILTO=rhead
# make beta in trackDb 3 x a week, prior to trackDbGlobal runs (on hgwdev)
5 5 * * 1,3,5 cd /cluster/home/qateam/trackDb; make beta > /dev/null
