# src/hg/makeDb/genbank/etc/build.crontab
#
# Crontab to run genbank nightly download, process, and align tasks.
#

SHELL=/bin/sh
MAILTO=markd@cse.ucsc.edu

# build genbank 5 minutes after midnight
5 0 * * *	/cluster/store5/genbank/etc/build-genbank


# try MGC status info daily at 3am.  Updated Wednesday 2 am EST, we check
# every day in case of failures.
0 3 * * *	/cluster/store5/genbank/etc/build-mgc
