# src/hg/makeDb/genbank/etc/dbload.crontab
#
# load database
#

SHELL=/bin/sh
MAILTO=markd@cse.ucsc.edu

# poll every hour at half-past
30 0-23 * * *	/scratch/genbank/etc/genbank-local-load




