# src/hg/makeDb/genbank/etc/gbdb-copy.crontab
#
# update data files from the build server to the gbdb server
#

SHELL=/bin/sh
MAILTO=markd@cse.ucsc.edu

# poll every hour at 50 minutes after
50 0-23 * * *	/scratch/genbank/etc/genbank-gbdb-copy




