# src/hg/makeDb/genbank/etc/hgdownload.crontab
#
# update download files on hgdownload
#

SHELL=/bin/sh
MAILTO=genbank@cse.ucsc.edu
PATH=/usr/local/bin:/usr/bin:/bin

# poll every hour at 5 minutes after
5 0-23 * * *	/export/genbank/etc/hgdownload-update
