# src/hg/makeDb/genbank/etc/hgdownload.crontab
#
# update download files on hgdownload
#

SHELL=/bin/sh
MAILTO=genbank@cse.ucsc.edu
PATH=/usr/local/bin:/usr/bin:/bin

# poll every 12 hours; a lot of time is allowed, otherwise the program
# will download as the genomes complete, generating lots of e-mail messages.
5 0,11,23 * * *	/genbank/etc/hgdownload-update
