# src/hg/makeDb/genbank/etc/hgdownload.crontab
#
# generate download files on hgdownload
# - this file need to work on both linux and solaris
#

# copy new download files from genbank server
# try hourly 
30 0-23 * * * /genbank/etc/hgdownload-copy
