# src/hg/makeDb/genbank/etc/hgdownload.crontab
#
# generate download files on hgdownload
#

SHELL=/bin/sh
MAILTO=genbank@cse.ucsc.edu
PATH=/usr/local/bin:/usr/bin:/bin

# copy new download files from genbank server
# try hourly 
30 */6 * * * /genbank/etc/hgdownload-copy
