# src/hg/makeDb/genbank/etc/hgwbeta.crontab
#
# dbload of hgwbeta
#

SHELL=/bin/sh
MAILTO=genbank@cse.ucsc.edu

# during the week only load mRNAs
# poll every 15 minutes 1-7
0,15,30,45 1-7 * * 1-5	/genbank/etc/hgwbeta-dbload

# on weekends, mRNAs and ESTS, and extFile table update
# poll every 15 minutes sat/sun
0,15,30,45 0-23 0-23 * * 0,6	/genbank/etc/hgwbeta-dbload -inclEsts
# -extFileUpdate


