# src/hg/makeDb/genbank/etc/hgwbeta.crontab
#
# dbload of hgwbeta
#

SHELL=/bin/sh
MAILTO=genbank@cse.ucsc.edu

# during the week only load mRNAs
# poll hourly 1-7
30 1-7 * * 1-5	/genbank/etc/hgwbeta-dbload

# on weekends, mRNAs and ESTS, and extFile table update
# poll hourly sat 0:30 - 23:30, sun 0:30 - 18:00
30 0-23 * * 0,6	/genbank/etc/hgwbeta-dbload -inclEsts
# -extFileUpdate


