# src/hg/makeDb/genbank/etc/hgwbeta.crontab
#
# dbload of hgwbeta
#

SHELL=/bin/sh
MAILTO=genbank@cse.ucsc.edu

# sunday-thrusday load mRNAs
# poll every 15 minutes 1-7
0,15,30,45 1-7 * * 0-4	/genbank/etc/hgwbeta-dbload

# friday night and saturday, both mRNAs and ESTS
# poll every 15 minutes fri/sat
0,15,30,45 0-14 * * 5,6	/genbank/etc/hgwbeta-dbload -inclEsts
# -extFileUpdate


