# src/hg/makeDb/genbank/etc/hgwdev.crontab
#
# load database on hgwdev server
#

SHELL=/bin/sh
MAILTO=genbank@cse.ucsc.edu

# mon-fri AM, load mRNAs
# poll every 15 minutes 1-6
0,15,30,45 1-6 * * 0-5	/genbank/etc/hgwbeta-dbload

# sat,sun AM both mRNAs and ESTS
# poll every 15 minutes midnight to 2pm
0,15,30,45 0-14 * * 6,7	/genbank/etc/hgwbeta-dbload -inclEsts
# -extFileUpdate


