# src/hg/makeDb/genbank/etc/hgwdev.crontab
#
# load database on hgwdev server
#

SHELL=/bin/sh
MAILTO=genbank@cse.ucsc.edu

# during the week only load mRNAs
# poll every 15 minutes 1-7
0,15,30,45 1-7 * * 1-5	/cluster/data/genbank/etc/hgwdev-dbload

# on weekends, both mRNAs and ESTS
# poll every 15 minutes sat/sun
0,15,30,45 0-23 * * 0,6	/cluster/data/genbank/etc/hgwdev-dbload -inclEsts
# -extFileUpdate


