# src/hg/makeDb/genbank/etc/mgc-dbload.crontab
#
# load database on MGC server
#

SHELL=/bin/sh
MAILTO=markd@cse.ucsc.edu

# poll every hour at 50 minutes after
50 0-23 * * *	/scratch/genbank/etc/genbank-gbdb-copy

# poll every hour at half-past
30 0-23 * * *	/scratch/genbank/etc/genbank-mgc-load




