# src/hg/makeDb/genbank/etc/mgc.crontab
#
# load database on MGC server
#

SHELL=/bin/sh
MAILTO=genbank@cse.ucsc.edu

# during the week and weekends we do both mRNAs and ESTs.  Only on this
# machine, MGC needs current ESTs and this is also nice sanity check for other
# machines.  Only do extFileUpdate on weekends.
30 22-23,00-04 * * 1-5	/genbank/etc/mgc-dbload -inclEsts
30 22-23,00-04 * * 0,6	/genbank/etc/mgc-dbload -inclEsts

# -extFileUpdate

