# src/hg/makeDb/genbank/etc/mgc.crontab
#
# load database on MGC server
#

SHELL=/bin/sh
MAILTO=markd@cse.ucsc.edu

# poll every hour at half-past
30 0-23 * * *	/genbank/etc/genbank-mgc-load




