# src/hg/makeDb/genbank/etc/mgc.crontab
#
# load database on MGC server
#

SHELL=/bin/sh
MAILTO=genbank@cse.ucsc.edu

# during the week and weekends we do both mRNAs and ESTs.  Only on this
# machine, MGC needs current ESTs and this is also nice sanity check for other
# machines.
# daily, 22:00-4:00
30 22-04 * * *	/genbank/etc/mgc-dbload -inclEsts

