# src/hg/makeDb/genbank/etc/rr.crontab
#
# dbload of round-robin servers.
#
SHELL=/bin/sh
MAILTO=genbank@cse.ucsc.edu

# mon-thru late afternoon, load mRNAs
# poll hourly 16:30 - 21:30, which is low usage period
30 16-21 * * 1-4	/genbank/etc/rr-dbload

# fri late afternoon on, load mRNAs and ESTS
30 16-23 * * 5	/genbank/etc/rr-dbload -inclEsts

# sat, all day load mRNAs and ESTS
# poll every hour
30 0-23 * * 6	/genbank/etc/rr-dbload -inclEsts

