# src/hg/makeDb/genbank/etc/rr.crontab
#
# dbload of round-robin servers.
#
SHELL=/bin/sh
MAILTO=genbank@cse.ucsc.edu

# sunday-thrusday load mRNAs
# poll hourly 16:30 - 21:30
30 16-21 * * 0-4	/genbank/etc/rr-dbload

# friday night and saturday, both mRNAs and ESTS
# poll hourly fri 16:30-23:30, sat all day
30 16-23 * * 5	/genbank/etc/rr-dbload -inclEsts
30 0-23 * * 6	/genbank/etc/rr-dbload -inclEsts
# -extFileUpdate

