# src/browserbox/var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.z5lDOp/crontab installed on Thu Mar 27 16:48:06 2014)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
SHELL=/bin/bash
# If you need to measure trash more often than just once a day,
# enable this measureTrash.sh script.  The example here runs every four hours
# on the even hours at 33 minutes after the hour

33 0-23/4 * * * /root/browserLogs/scripts/measureTrash.sh /root/browserLogs/scripts/virtualBox.txt

# if you need to clean trash more often than just once a day,
# enable this trashCleanMonitor.sh script.  The example here runs every
#  four hours on the even hours at 55 minutes after the hour:

55 1-23/4 * * * /root/browserLogs/scripts/trashCleanMonitor.sh /root/browserLogs/scripts/virtualBox.txt searchAndDestroy

# These scripts can lock each other out to prevent overload protection
# on the filesystem where trash is located.  The measureTrash.sh script
# doesn't mind being locked out by the trash cleaner, it will just skip
# the measurement without error.
# But the trash cleaner does not like being locked out by
# the measureTrash script and will be considered an error.
# On really busy heavily used systems, beware of their run times overlapping
# each other.

# update all browser CGIs, gbdb files and mysql tables, every four hours, at x:30
# http://stackoverflow.com/questions/9049460/cron-jobs-and-random-times-within-giving-hours
30 1-23/4 * * * /root/cronUpdate.sh

# clear tmp once per day, to make sure that hgMirror does not have too old
# cached copies
0 3 * * * /root/cronCleanTmp.sh

