# src/browserbox/var/spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XCE7ef/crontab installed on Thu Oct  9 12:12:16 2014)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
SHELL=/bin/bash
# If you need to measure trash more often than just once a day,
# enable this measureTrash.sh script.  The example here runs every four hours
# on the even hours at 33 minutes after the hour
# added sleep statement so it doesn't start right away when the box is booting

33 0-23/4 * * * sleep 240 ; /root/browserLogs/scripts/measureTrash.sh /root/browserLogs/scripts/virtualBox.txt

# if you need to clean trash more often than just once a day,
# enable this trashCleanMonitor.sh script.  The example here runs every
#  four hours on the even hours at 55 minutes after the hour:

55 1-23/4 * * * sleep 240 ; /root/browserLogs/scripts/trashCleanMonitor.sh /root/browserLogs/scripts/virtualBox.txt searchAndDestroy

# These scripts can lock each other out to prevent overload protection
# on the filesystem where trash is located.  The measureTrash.sh script
# doesn't mind being locked out by the trash cleaner, it will just skip
# the measurement without error.
# But the trash cleaner does not like being locked out by
# the measureTrash script and will be considered an error.
# On really busy heavily used systems, beware of their run times overlapping
# each other.

# update all browser CGIs, gbdb files and mysql tables, every minute
# http://stackoverflow.com/questions/9049460/cron-jobs-and-random-times-within-giving-hours
*/4 * * * * /root/cronUpdate.sh

# clear tmp every night, to make sure that hgMirror does not have too old
# cached copies
0 3 * * * /root/cronCleanTmp.sh

# need to make sure there always is an update script 
*/3 * * * * if [[ ! -s /root/updateBrowser.sh ]] ; then wget http://hgdownload.gi.ucsc.edu/gbib/updateBrowser.sh -O /root/updateBrowser.sh ; chmod a+x /root/updateBrowser.sh; echo redownloaded updateBrowser.sh ;fi

# we check only once if UDR is available, as it can take a while to time out
@reboot sleep 60; udr rsync -avP hgdownload.gi.ucsc.edu::gbib/lastUpdate /tmp/useUdr > /dev/null 2> /dev/null
